# .devcontainer/supervisor.d/opencode-supervisor.conf
[program:opencode]
command=/usr/local/bin/mise exec bun -- bunx opencode-ai serve --hostname 0.0.0.0
directory=/workspaces/agentic-plugins
autostart=true
startsecs=5
autorestart=true
startretries=3
stderr_logfile=/var/log/opencode.err.log
stdout_logfile=/var/log/opencode.log
user=vscode
environment=PATH="/home/vscode/.local/bin:/usr/local/bin:/usr/bin:/bin",HOME="/home/vscode"